# SYSTEMD/.config/systemd/user/backup.service
[Unit]
Description=Backup
RequiresMountsFor=/mnt/Data
Wants=backup.timer

[Service]
Type=simple
ExecStart=%h/.local/bin/backup.sh

